# scripts/99-guitar.rules
SUBSYSTEMS=="usb", ATTRS{idVendor}=="08bb", ATTRS{idProduct}=="2902", RUN+="/bin/sh -c '/bin/echo -Chw:CODEC,0 -Phw:CODEC > /tmp/config'"
