# scripts/99-midi_keyboard.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="1acc", ATTRS{idProduct}=="1633", ACTION=="add", RUN+="/bin/sh -c 'echo > /tmp/midi'"
